# docker/quadlets/redis.container
# Podman Quadlet: Redis for Swedeb Celery task broker
# Deploy to ~/.config/containers/systemd/ (rootless) or /etc/containers/systemd/ (root)

[Unit]
Description=Redis broker for Swedeb Celery workers
Requires=swedeb-staging-app.network
After=swedeb-staging-app.network

[Container]
Image=docker.io/redis:7-alpine
ContainerName=swedeb-redis-staging
Network=swedeb-staging-app.network
Volume=redis_data/redis:/data

# Persist data across restarts
Exec=redis-server --appendonly yes

[Service]
Restart=always
TimeoutStartSec=30

[Install]
WantedBy=default.target
